# services/simulator/deploy/taqinor-simulator.service
# TAQINOR Solar Quote Simulator — systemd unit (Hetzner production).
# Installed at /etc/systemd/system/taqinor-simulator.service
#
# The app is a FastAPI (ASGI) app run with uvicorn. It binds to the docker
# bridge gateway 172.18.0.1 (the host, as seen from the Caddy container on the
# erp-agentique_default network). External traffic on eth0 cannot reach that
# address, so the app is NOT publicly exposed — only Caddy reaches it, and Caddy
# terminates TLS. Restart=always recovers from the brief boot race where the
# docker bridge is not yet up when this service first starts.
[Unit]
Description=TAQINOR Solar Quote Simulator (FastAPI/uvicorn)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=simulator
Group=simulator
WorkingDirectory=/opt/taqinor-simulator
Environment=PYTHONUNBUFFERED=1
# Server-only secret (SIMULATOR_SECRET_KEY) — git-ignored, chmod 600, owned by
# the simulator user. Required (no leading '-'): the service fails to start if
# it is missing, so the app can never run with a missing/blank signing secret.
EnvironmentFile=/opt/taqinor-simulator/secret.env
ExecStart=/opt/taqinor-simulator/venv/bin/uvicorn main:app --host 172.18.0.1 --port 8090 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=3

# Hardening — the app only needs to read/write its own data dir.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/taqinor-simulator

[Install]
WantedBy=multi-user.target
